# smart-mirror.service
[Unit]
Description=Smart Mirror startup service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/home/pi/smart-mirror.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
